# deploy/leadforce.service
[Unit]
Description=LeadForce Flask Service
After=network.target

[Service]
User=root
WorkingDirectory=/root/LeadForcePython
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m gunicorn -w 1 -b 0.0.0.0:12345 --timeout 120 app:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
